# init/ivc-server.service
[Unit]
Description=BMS IVC Server
Requires=docker.service

[Service]
Type=notify
RemainAfterExit=yes

ExecStartPre=-/usr/bin/docker create --name ivc-server --restart on-failure:5 -v /opt/ivc:/opt/ivc -p 9987-9989:9987-9989/udp 108vfs/ivc
ExecStart=/usr/bin/docker start ivc-server

ExecStop=/usr/bin/docker stop ivc-server
ExecStopPost=/usr/bin/docker rm -f ivc-server

ExecReload=/usr/bin/docker restart ivc-server

SuccessExitStatus=0 143

[Install]
WantedBy=multi-user.target
Alias=ivc.service

